# debian/prerm
#!/bin/sh
invoke-rc.d chilli stop

if [ "$1" = purge ]; then
    update-rc.d chilli remove
fi
